# DDE_UAT/2026-01-29/files/etc/profile.d/vc_notes.sh
#
#	Server notes
#
#	Alan Lodewyk
#	2017-10-13
#
#	This script displays the /root/notes.txt
#
gidName=`id -gn`
if [ "$gidName" = "root" ] ; then
	[ -f /root/notes.txt ] && cat /root/notes.txt
fi
